# bot/crontab
0 */4 * * * root /app/run.sh | sed -u 's/\%/\%\%/g' | xargs -d '\n' -n 1 -I {} date +"[\%Y-\%m-\%dT\%H:\%M:\%S.\%6N\%:z] {}" >> /var/log/cron.log 2>&1
